# systemd/user/openviking-wechat-archive-server.service
[Unit]
Description=OpenViking WeChat Archive HTTP Service
After=network-online.target openviking-local-embed.service openviking-local-rerank.service
Wants=network-online.target openviking-local-embed.service openviking-local-rerank.service

[Service]
Type=simple
WorkingDirectory=/home/nx/github/OpenViking
ExecStart=/home/nx/github/OpenViking/scripts/run_wechat_archive_http_service.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
